# docker/rootfs/etc/s6-overlay/s6-rc.d/leapmux/run
#!/command/execlineb -P
importas MODE LEAPMUX_MODE
/usr/local/bin/leapmux ${MODE} -data-dir /data
