# docker/rootfs/etc/s6-overlay/s6-rc.d/leapmux/run
#!/command/execlineb -P
/usr/local/bin/leapmux -addr :4327 -data-dir /data
